# roles/systemd/templates/wpbootstrap.service.j2
[Unit]
Description=wpbootstrap
After=cloud-init.target

[Service]
Type=oneshot
ExecStart=/usr/bin/nohup /usr/bin/wpbootstrap.sh
PIDFile=/var/run/wpbootstrap.pid

[Install]
WantedBy=multi-user.target
